# rootfs/etc/services.d/01rpcbind/run
#!/usr/bin/execlineb -P
/sbin/rpcbind -d
