# rootfs/etc/service.d/rpc.nfsd/run
#!/usr/bin/execlineb -P
/usr/sbin/rpc.nfsd 8

