# rootfs/etc/service.d/rpc.mountd/run
#!/usr/bin/execlineb -P
# support only nfs version 3 (no 2 no 4)
/usr/sbin/rpc.mountd -V 3 -N 2 -N 4
