# rootfs/etc/services.d/02rpc.statd/run
#!/usr/bin/execlineb -P
/usr/sbin/rpc.statd -F -d -L
