# rootfs/etc/services.d/03rpc.nfsd/finish
#!/usr/bin/execlineb -P
/usr/sbin/rpc.nfsd 0 
echo 'rpc.nfsd finish called!'

