# rootfs/etc/service.d/rpc.statd/run
#!/usr/bin/execlineb -P
/usr/sbin/rpc.statd --no-notify
