# rootfs/etc/services.d/04rpc.mountd/run
#!/usr/bin/execlineb -P
# support only nfs version 3 (no 2 no 4)
/usr/sbin/rpc.mountd -V3 -N2 -N4 -F
